# longevity-backend.service
[Unit]
Description=Longevity OS — FastAPI backend
After=network.target

[Service]
Type=simple
User=dcb
WorkingDirectory=/home/dcb/longevity
EnvironmentFile=/home/dcb/longevity/.env
ExecStart=/home/dcb/.local/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
